# rocky_sim/rootfs/etc/s6-overlay/s6-rc.d/mariadb/run
#!/command/with-contenv sh
set -eu

exec /usr/libexec/mariadbd --user=mysql --console
