# rootfs/etc/s6-overlay/s6-rc.d/munge/run
#!/command/with-contenv sh
set -eu

exec 2>&1
exec /usr/sbin/munged -F -f
