# rootfs/etc/s6-overlay/s6-rc.d/slurmd/run
#!/command/with-contenv sh
set -eu

exec 2>&1
until [ -S /run/munge/munge.socket.2 ]; do
    sleep 1
done

exec /usr/sbin/slurmd -D -vv -f /etc/slurm/slurm.conf
