# rocky_sim/rootfs/etc/s6-overlay/s6-rc.d/slurmd/run
#!/command/with-contenv sh
set -eu

until [ -S /run/munge/munge.socket.2 ]; do
    sleep 1
done

exec /usr/sbin/slurmd -D -vv
