# rootfs/etc/s6-overlay/s6-rc.d/mariadb-init/up
/etc/s6-overlay/scripts/mariadb-init.sh
